# arch-armhf/DEBIAN/prerm
#!/bin/sh

set -e

daemon=tiny-ssdp

invoke-rc.d "$daemon" stop || :

